# base_python/bootstrap/supervisor-bootstrap.conf
[program:filebeat]
priority = 1
command = filebeat -c /filebeat.yml