# app/supervisor-app.conf
[program:test]
command = sleep infinity

[program:testpy]
command = ./venv/bin/python example.py